# test/fail.service
[Unit]
Description=test service that writes to stderr

[Service]
Type=simple
ExecStart=/usr/local/bin/fail
User=nobody
