# journal-monitor.timer
[Unit]
Description=run journal-monitor periodically

[Timer]
OnBootSec=15min
OnActiveSec=15min

[Install]
WantedBy=timers.target
